# crontab
0 * * * * /.../iq/manage.py checkaccount
